# rootfs/etc/s6-overlay/s6-rc.d/collab/run
#!/command/with-contenv sh
cd /app
exec ./node_modules/.bin/tsx collab/server.ts
